# confs/nest.conf
[program:nest]
command = /talkativot/confs/nest_runner.sh
directory=/talkativot
autostart=true
autorestart=true
startsecs=10
stopwaitsecs = 600
killasgroup=true
priority=1000
